# crontab
# 2017年 6月12日 星期一 18时30分42秒 CST
SHELL=/usr/local/bin/bash
30 10 * * * /usr/local/bin/brew update &>> /tmp/cron.log &
